# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys
Icon=
Name=xbindkeys
Path=
Terminal=True
Type=Application
